# v2ray-customized.service
[Unit]
Description=The V2Ray Client Service
Documentation=https://www.v2fly.org/
After=network.target nss-lookup.target

[Service]
Type=simple
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ExecStartPre=/usr/bin/v2ray -c /etc/v2ray/MY_V2RAY_CONFIG.json -test
ExecStart=/usr/bin/v2ray -c /etc/v2ray/MY_V2RAY_CONFIG.json
Restart=on-failure
RestartPreventExitStatus=23

[Install]
WantedBy=multi-user.target